# systemd/component-storage.service
[Unit]
Description=Component Storage - Electronic Component Inventory System
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/component-storage/backend
Environment=NODE_ENV=production
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=component-storage

[Install]
WantedBy=multi-user.target
