# systemd/component-storage.service
[Unit]
Description=Component Storage FastAPI
After=network.target

[Service]
User=component-storage
Group=component-storage
WorkingDirectory=/opt/component-storage/current/backend/app
EnvironmentFile=/etc/component-storage/.env
ExecStart=/opt/component-storage/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
